# debian/postinst
#!/bin/sh
update-rc.d behaviours start 20 2 3 4 5 . stop 1 6 .
